# Change Default Password/com.admin.checkdefaultpassword.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this LaunchAgent -->
    <key>Label</key>
    <string>com.admin.checkdefaultpassword</string>

    <!-- Script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/bin/check_default_password.sh</string>
    </array>

    <!-- Run every time the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Run again on subsequent logins (not just once ever) -->
    <key>LaunchOnlyOnce</key>
    <false/>

    <!-- Delay a few seconds to let the desktop fully load before showing dialogs -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
